# deploy/teatree-watchdog.timer
# Fires the teatree self-heal watchdog every few minutes. Because the timer is a
# host-level systemd unit — OUTSIDE the compose stack it watches — a full stack
# outage (the init crash-loop that froze the factory) cannot take the watchdog
# down with it.
[Unit]
Description=Run the teatree self-heal watchdog every few minutes
Documentation=https://github.com/souliane/teatree

[Timer]
# First pass shortly after boot, then every 5 minutes while the box is up.
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
