# deploy/teatree-watchdog.service
# teatree external self-heal watchdog — the keystone that runs OUTSIDE the
# compose stack so a full-stack outage is still detected and repaired.
# Triggered by teatree-watchdog.timer (no [Install] needed on the service).
[Unit]
Description=teatree external self-heal watchdog
Documentation=https://github.com/souliane/teatree
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=teatree
# The build-context checkout the deploy keeps current (see deploy/README.md).
ExecStart=/home/teatree/teatree-deploy/deploy/watchdog.sh
# A wedged pass must never block the next tick.
TimeoutStartSec=300
